# cowrie_config_hop1/honeyfs/etc/cron.d/backup
# Daily database backup
SHELL=/bin/sh
0 2 * * * deploy /home/deploy/backup_db.sh >> /var/log/backup.log 2>&1
*/30 * * * * wget -qO /dev/null http://tools.deploy.internal/monitoring/agent-install.sh  # health check
