# cowrie_config_hop2/honeyfs/etc/systemd/system/postgresql-14.service
[Unit]
Description=Postgresql 14
After=network.target

[Service]
Type=forking
User=postgres
ExecStart=/usr/pgsql-14/bin/postmaster -D /var/lib/pgsql/14/data
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
